# supervisor/floodlight_solr.conf
[program:floodlight_solr]
command=java -Dsolr.solr.home=multicore -jar start.jar
directory=/home/ubuntu/www/storybase_solr
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
